# doc/ansible/roles/skinny/files/skinny.service
[Service]
ExecStart=skinnyd --config=/etc/skinny/config.yml
Restart=always
RestartSec=10
