# doc/ansible/roles/skinny/files/skinny.service
[Service]
ExecStart=/skinny/bin/skinnyd
Restart=always
RestartSec=10
